# units/fixinputs.path
# Inspired by http://jasonwryan.com/blog/2014/01/20/udev/
# and http://www.ocsmag.com/2015/09/02/monitoring-file-access-for-dummies/
[Unit]
Description=Triggers the service that sets up external keyboard

[Path]
PathChanged=/dev/input/

[Install]
WantedBy=multi-user.target
