# etc/systemd/system/systemd-logind.service.d/override.conf
[Service]
# Allow symlinks into the home directory.
ProtectHome=no
